# app-admin/cobbler/files/cobblerd
#!/sbin/runscript
#
# cobblerd    Cobbler helper daemon
###################################

#
# processname: /usr/bin/cobblerd

depend() {
	use logger apache2
        need net
}


start() {
	CONFIG_ARGS="-F"
	ebegin "Starting cobblerd"
        start-stop-daemon --start --quiet --pidfile /run/cobblerd.pid -b -m \
		/usr/bin/cobblerd -- ${CONFIG_ARGS}
        eend $?
}

stop() {
	ebegin "Stopping cobblerd"
        start-stop-daemon --stop --quiet --pidfile /run/cobblerd.pid  
        eend $?
}
